# private_dot_config/systemd/user/do_sync.service
[Unit]
Description=Sync selected files

[Service]
Type=simple
ExecStart=%h/bin/do_sync
