# private_dot_config/systemd/user/do_sync.timer
[Unit]
Description=Sync selected files

[Timer]
OnBootSec=15min
OnUnitInactiveSec=1h
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
