# private_dot_config/systemd/user/gnome-polkit.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
ConditionEnvironment=|XDG_CURRENT_DESKTOP=Hyprland
ConditionEnvironment=|XDG_CURRENT_DESKTOP=niri

[Service]
ExecStart=/usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
Restart=on-failure

[Install]
WantedBy=graphical-session.target
