# private_dot_config/systemd/user/wlsunset.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
ConditionEnvironment=|XDG_CURRENT_DESKTOP=niri
ConditionEnvironment=|XDG_CURRENT_DESKTOP=Hyprland

[Service]
ExecStart=/usr/bin/wlsunset -l 36.3 -L 116.3 -t 4000
Restart=always

[Install]
WantedBy=graphical-session.target
